# task05/10-local.rules
SUBSYSTEMS=="usb", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0745", RUN+="modprobe task05"
